# .infrastructure/crontab
*/5 * * * * su www-data -s /bin/bash -c "cd /var/www/html && php ./bin/console background-process"
