# .infrastructure/etc/supervisor/conf.d/logs.conf
[program:synchronous-background-process]
command=tail -f /var/log/uptime-admin-board.log
autorestart=true
autostart=true
user=www-data
redirect_stderr=true
redirect_stdout=true
