# .infrastructure/etc/supervisor/conf.d/synchronous-background-process.conf
[program:synchronous-background-process]
command=bash /synchronous-background-process.sh
autorestart=false
autostart=true
user=www-data
redirect_stderr=true
redirect_stdout=true
